# systemd/otp-monitor.service
[Unit]
Description=OTP Relay Monitor (phone watcher + alert forwarder)
After=network-online.target otp-relay.service
Wants=network-online.target otp-relay.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/otp-relay
EnvironmentFile=/opt/otp-relay/.env
ExecStart=/opt/otp-relay/venv/bin/python3 monitor.py
Restart=always
RestartSec=10
PrivateTmp=true
ReadWritePaths=/opt/otp-relay/data

[Install]
WantedBy=multi-user.target
